# justlayme-tunnel.service
[Unit]
Description=JustLayMe Cloudflare Tunnel
After=network-online.target justlayme-api.service
Wants=network-online.target
Requires=justlayme-api.service
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=15
User=fastl
Group=fastl
WorkingDirectory=/home/fastl/JustLayme
ExecStart=/usr/local/bin/cloudflared tunnel run --url http://localhost:3000 justlayme
StandardOutput=journal
StandardError=journal
SyslogIdentifier=justlayme-tunnel

# Restart policy for network issues
Restart=always
RestartSec=15
StartLimitBurst=10
StartLimitIntervalSec=120

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=false
ProtectSystem=strict
ReadWritePaths=/home/fastl

[Install]
WantedBy=multi-user.target